# controlpod.service
[Unit]
Description=KCLF Control Pod
After=network-online.target local-fs.target dev-ttyUSB0.device dev-i2c-1.device
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/ControlPod
ExecStart=/home/pi/control_pod/env/bin/python /home/pi/control_pod/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment=HEARTBEAT_PATH=/run/controlpod/heartbeat
Environment=LAST_SEND_PATH=/run/controlpod/last_send
Environment=SHUTDOWN_PATH=/run/controlpod/shutdown
RuntimeDirectory=controlpod
RuntimeDirectoryMode=0755
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target

