# systemd/controlpod-healthcheck.timer
[Unit]
Description=Run ControlPod heartbeat health check every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=1min

[Install]
WantedBy=timers.target
